# templates/agent/systemd.puppet-run.service.erb
#
# Managed by Puppet
#
[Unit]
Description=Systemd Timer Service for Puppet Agent
Wants=network.target network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
ExecStart=<%= @command %>
SuccessExitStatus=2
User=root
Group=root
