# templates/agent/systemd.puppetcron.service.erb
[Unit]
Description=Systemd Timer Service for Puppet Agent

[Service]
Type=oneshot
ExecStart=/usr/bin/puppet agent --no-daemonize --onetime
User=root
Group=root
