# templates/agent/systemd.puppetcron.timer.erb
<% times[0] = times[0].join(",") if times[0].kind_of? Array %>
<% times[1] = times[1].join(",") if times[1].kind_of? Array %>
[Unit]
Description=Systemd Timer for Puppet Agent

[Timer]
OnActiveSec=0m
OnCalendar=<%= times[0] %>:<%= times[1] %>
Persistent=true

[Install]
WantedBy=timers.target
