# templates/agent/systemd.puppetcron.service.erb
[Unit]
Description=Systemd Timer Service for Puppet Agent

[Service]
Type=oneshot
ExecStart=<%= @command %>
User=root
Group=root
