# templates/agent/systemd.puppet-run.timer.erb
#
# Managed by Puppet
#
[Unit]
Description=Systemd Timer for Puppet Agent

[Timer]
OnCalendar=*-*-* <%= Array(@_hour).join(',') %>:<%= Array(@_minute).join(',') %>:00 <%= @timezone %>
Persistent=true
RandomizedDelaySec=<%= @randomizeddelaysec %>

[Install]
WantedBy=timers.target
