# templates/agent/systemd.puppet-run.service.erb
#
# Managed by Puppet
#
[Unit]
Description=Systemd Timer Service for Puppet Agent

[Service]
Type=oneshot
ExecStart=<%= @command %>
SuccessExitStatus=2
User=root
Group=root
