# templates/agent/systemd.puppet-run.timer.erb
[Unit]
Description=Systemd Timer for Puppet Agent

[Timer]
OnCalendar=*-*-* <%= Array(@times[0]).join(',') %>:<%= Array(@times[1]).join(',') %>:00
Persistent=true

[Install]
WantedBy=timers.target
