# templates/agent/systemd.puppet-run.timer.erb
#
# Managed by Puppet
#
[Unit]
Description=Systemd Timer for Puppet Agent

[Timer]
OnCalendar=*-*-* <%= Array(@times[0]).join(',') %>:<%= Array(@times[1]).join(',') %>:00
Persistent=true
RandomizedDelaySec=<%= @randomizeddelaysec %>

[Install]
WantedBy=timers.target
