# templates/systemd.service.erb
[Unit]
Description=Run Puppet Agent --onetime

[Service]
ExecStart=<%= @command %>
SuccessExitStatus=2
Type=oneshot

